# bios-bmc-smm-error-logging/service_files/com.google.gbmc.bmcbiosbuffer.service.in
[Unit]
Description=BMC BIOS shared buffer daemon

[Service]
Type=simple
ExecStart=@BINDIR@/bmcbiosbuffer

[Install]
WantedBy=multi-user.target